# .crontab
# custom cron job
* * * * *  /usr/bin/php /var/www/bin/console app:cash-sync
